# Site Web/crontab
* * * * * /var/www/html/action.php?arrosage&case=5